# deploy/systemd/arcsent-watchdog.timer
[Unit]
Description=Run ArCsent Watchdog

[Timer]
OnBootSec=2m
OnUnitActiveSec=2m
Unit=arcsent-watchdog.service

[Install]
WantedBy=timers.target
